# scripts/systemd/celery-short.service
[Unit]
Description=Celery Worker (Short-lived Queues)
After=network.target redis.service

[Service]
Type=forking
User=training
Group=training
WorkingDirectory=/opt/training_platform

# Environment
Environment=REDIS_URL=redis://localhost:6381/0
Environment=DATABASE_URL=postgresql://training_platform:password@localhost:5432/training_platform
Environment=PYTHONPATH=/opt/training_platform

# Command
ExecStart=/opt/training_platform/venv/bin/celery \
    -A training_platform.core.celery_config worker \
    -Q default,evaluation,preprocessing,maintenance \
    -c 4 \
    -l info \
    -n short@%h \
    --pidfile=/var/run/celery/short.pid \
    --logfile=/var/log/celery/short.log

# PID file
PIDFile=/var/run/celery/short.pid

# Restart policy
Restart=always
RestartSec=10

# Resource limits (optional)
# LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
